# wifi/configs/nlocker.service
[Unit]
Description=NLocker

[Service]
ExecStart=#PATH#/nlocker/main.py

[Install]
WantedBy=multi-user.target
