# wifi/configs/raspberry_wifi_api.service
[Unit]
Description=Raspberry WiFi API

[Service]
ExecStart=#PATH#/api/main.py

[Install]
WantedBy=multi-user.target
